# resources/scripts-deb/prerm
#!/bin/bash
/bin/systemctl stop stalwart-jmap
/bin/systemctl disable stalwart-jmap
deluser stalwart-jmap
